# template/crontab
0 16 * * * /usr/local/duplic/dupl.sh usb
0 18 * * * /usr/local/duplic/dupl.sh usb1

