# template/crontab
0 16 * * * /usr/local/duplic/dupl.sh usb >>/usr/local/duplic/log/cron.out 2>&1
0 18 * * * /usr/local/duplic/dupl.sh usb1 >>/usr/local/duplic/log/cron.out 2>&1

